# hdp_zookeeper_node/conf/supervisor/zookeeper-supervisord.conf
[program:zookeeper]
command=/usr/lib/zookeeper/bin/zkServer.sh start-foreground 
autorestart=true
startretries=50000
log_stderr=true
logfile=/var/log/zookeeper/zookeeper.log
user=zookeeper

