# hdp_hive_hbase_node/conf/supervisor/hregionserver-supervisord.conf
[program:hregionserver]
command=hbase regionserver start
autorestart=true
startretries=50000
log_stderr=true
logfile=/var/log/hbase/hbase.log
user=hbase

